# sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
autorestart=true
priority=10
startsecs=2
startretries=5

[program:fluxbox]
command=bash -c "sleep 2 && fluxbox"
environment=DISPLAY=":99"
autorestart=true
priority=20
startsecs=3
startretries=5

[program:x11vnc]
command=bash -c "sleep 3 && x11vnc -display :99 -forever -nopw -shared -rfbport 5900 -xkb"
autorestart=true
priority=30
startsecs=4
startretries=5

[program:websockify]
command=bash -c "sleep 4 && websockify --web /usr/share/novnc 6080 localhost:5900"
autorestart=true
priority=40
startsecs=5
startretries=5
